# data_record.conf
[program:record]
directory=/home/ubuntu/Grafana_DashBoard_v2
autostart=true
autorestart=true
startretries=10
user=ubuntu
logfile_maxbytes=0
logfile_backups=0
command=python3 app.py
